# msi-ec-custom.service
[Unit]
Description=Load custom MSI-EC module for fan control
After=multi-user.target
Before=graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Remover m칩dulo del kernel si est치 cargado
ExecStartPre=-/sbin/modprobe -r msi-ec
ExecStartPre=-/sbin/rmmod msi_ec
# Cargar m칩dulo personalizado
ExecStart=/sbin/insmod /home/manusoftar/Git/msi-ec/msi-ec.ko
# Descargar m칩dulo al apagar
ExecStop=/sbin/rmmod msi_ec

[Install]
WantedBy=multi-user.target
